# deploy/com.voitures-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.voitures-bot</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/alex/CascadeProjects/VOITURES/voitures-bot/venv/bin/python</string>
        <string>scripts/run_prod.py</string>
        <string>--loop</string>
        <string>--interval</string>
        <string>15</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/alex/CascadeProjects/VOITURES/voitures-bot</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/alex/CascadeProjects/VOITURES/voitures-bot/logs/bot.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/alex/CascadeProjects/VOITURES/voitures-bot/logs/bot.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alex/CascadeProjects/VOITURES/voitures-bot/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
